# fedora-rewrite/skel/etc/NetworkManager/dispatcher.d/90-start-salt
#!/usr/bin/env python3

import subprocess
import sys


if len(sys.argv) < 3 or sys.argv[2] != "up":
    exit(0)

subprocess.check_call(["systemctl", "start", "salt-minion"])