# config/includes.chroot/etc/init.d/polyprog-network
#!/bin/bash

### BEGIN INIT INFO
# Provides:          polyprog-network
# Required-Start:    $network $named
# Required-Stop:     $network $named
# Default-Start:     S
# Default-Stop:      
# Short-Description: Connects to the network.
# Description:       Connects to the network. Knows about EPFL configuration.
### END INIT INFO

. /lib/lsb/init-functions

N=/etc/init.d/polyprog-network

wait_for_connectivity() {
    log_action_begin_msg "Waiting for network connectivity"
    local n_tries=1
    local max_tries=3

    # has_connectivity is a return value; zero if connection succeeds
    ping -c1 polyprog.epfl.ch >/dev/null 2>&1
    local has_connectivity=$?

    while [ "$has_connectivity" -ne 0 ] && [ "$n_tries" -lt "$max_tries" ]; do
        log_action_cont_msg " waiting ($n_tries / $max_tries)"
        sleep 3
        ping -c1 polyprog.epfl.ch >/dev/null 2>&1
        has_connectivity=$?
        n_tries=$((n_tries + 1))
    done

    log_action_end_msg $has_connectivity
    return $has_connectivity
}

start_network() {
    if ! wait_for_connectivity; then
        dhclient -r wlan0
        pkill wpa_supplicant

        wpa_supplicant -i wlan0 -c /etc/wpa_supplicant/wpa_supplicant.conf -B
        sleep 8
        dhclient wlan0
    else
        return 0
    fi
}

stop_network() {
    dhclient -r wlan0
    pkill wpa_supplicant
}

case "$1" in
  start)
        log_action_begin_msg "Connecting to the network"
        start_network
        log_action_end_msg $?
        ;;
  stop)
        log_action_begin_msg "Disconnecting from the network"
        stop_network
        log_action_end_msg $?
        ;;
  reload|restart|force-reload)
        log_action_begin_msg "Re-connecting to the network"
        stop_network
        start_network
        log_action_end_msg $?
        ;;
  *)
        echo "Usage: $N start|stop" >&2
        exit 1
        ;;
esac

exit 0
